# src/systemd/.config/systemd/user/devilspie.service
[Unit]
Description=devilspie

[Service]
ExecStart=/usr/bin/devilspie -a

[Install]
WantedBy=graphical.target
WantedBy=default.target
