# src/systemd/.config/systemd/user/gitfiti.timer
[Unit]
Description=Refresh drawing on Github calendar weekly.

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
