# src/freedesktop/.config/autostart/xmodmap.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/xmodmap /home/eqyiel/.Xmodmap
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
Name=Xmodmap
